# gpscontrol4u.conf
[program:gpscontrol4u]
command=/home/systemd/gpscontrol4u/.venv/bin/python /home/systemd/gpscontrol4u/manage.py runserver 0.0.0.0:7008
directory=/home/systemd/gpscontrol4u
user=systemd
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/gpscontrol4u.err.log
stdout_logfile=/var/log/gpscontrol4u.out.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
stderr_logfile_backups=3
stdout_logfile_backups=3
environment=PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/home/systemd",USER="systemd",DJANGO_SETTINGS_MODULE="marketplace_backend.settings",DATABASE_HOST="localhost",DATABASE_PORT="3306",DATABASE_USER="gpscontrol",DATABASE_PASSWORD="qazwsxedc",DATABASE_NAME="gpscontrol4u",EMAIL_BACKEND="django.core.mail.backends.smtp.EmailBackend",EMAIL_HOST="pro.turbo-smtp.com",EMAIL_PORT="465",EMAIL_USE_TLS="False",EMAIL_USE_SSL="True",EMAIL_HOST_USER="datacollect@madd.com.mx",EMAIL_HOST_PASSWORD="GPSc0ntr0l01",DEFAULT_FROM_EMAIL="GPSControl4U <orders@gpscontrol.com.mx>",SECRET_KEY="your-super-secret-key-change-this-in-production",DEBUG="True",ALLOWED_HOSTS="localhost,127.0.0.1,mp.armaddia.lat",STRIPE_TEST_PUBLIC_KEY="pk_test_YOUR_ACTUAL_PUBLIC_KEY_HERE",STRIPE_TEST_SECRET_KEY="sk_test_YOUR_ACTUAL_SECRET_KEY_HERE",STRIPE_WEBHOOK_SECRET="whsec_YOUR_WEBHOOK_SECRET_HERE",MERCADO_PAGO_ACCESS_TOKEN="APP_USR-1747665656299106-112018-3c4e4831901495497d12b99bd8982350-184251616",MERCADO_PAGO_PUBLIC_KEY="APP_USR-70a80acf-4bd0-4ac3-b3e9-2f77456c1bb3",MERCADO_PAGO_SANDBOX="False",MERCADO_PAGO_TEST_BUYER_EMAIL="TESTUSER632906111@testuser.com",MERCADO_PAGO_TEST_BUYER_PASSWORD="PVOlNBKTEK",DOMAIN="store.gpscontrol4u.com",CLOUD_PLATFORM="RAILWAY",LANGUAGES="en,es",REDIS_URL="redis://localhost:6379/0",EXTERNAL_API_BASE_URL="https://api2ego.elisasoftware.com.mx",EXTERNAL_API_USERNAME="AdmGPScontrol4u",EXTERNAL_API_PASSWORD="GPSc0ntr0l4u*",EXTERNAL_API_STORE="GPScontrol4U"
